# .config/autostart/autostart-commands.desktop
[Desktop Entry]
Type=Application
Name=lucy-startup-commands
Exec=/home/lucy/.config/autostart-programs.sh
